# Trac 0.12-svn/supervisord.conf
[supervisord]
nodaemon=true

[program:svnserve]
command=/usr/bin/svnserve --daemon --foreground --root /svn

[program:tracd]
command=/usr/local/bin/tracd --single-env --port 80 /trac
